# deploy/supervisord.futurefunded.conf
[unix_http_server]
file=/home/elCUCO/futurefunded/.supervisor.sock

[supervisord]
logfile=/home/elCUCO/futurefunded/.logs/supervisord.log
logfile_maxbytes=20MB
logfile_backups=5
pidfile=/home/elCUCO/futurefunded/.supervisord.pid
childlogdir=/home/elCUCO/futurefunded/.logs
nodaemon=false

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///home/elCUCO/futurefunded/.supervisor.sock

[program:futurefunded-app]
directory=/home/elCUCO/futurefunded
command=/home/elCUCO/futurefunded/.venv/bin/gunicorn --workers 1 --threads 8 --bind 127.0.0.1:5000 --timeout 120 --graceful-timeout 30 --keep-alive 5 --worker-class gthread --access-logfile /home/elCUCO/futurefunded/.logs/gunicorn-access.log --error-logfile /home/elCUCO/futurefunded/.logs/gunicorn-error.log wsgi:app
user=elCUCO
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=20
stopsignal=TERM
environment=PYTHONUNBUFFERED="1"
stdout_logfile=/home/elCUCO/futurefunded/.logs/supervisor-app.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5
stderr_logfile=/home/elCUCO/futurefunded/.logs/supervisor-app.err.log
stderr_logfile_maxbytes=20MB
stderr_logfile_backups=5

[program:futurefunded-cloudflared]
command=/usr/bin/cloudflared tunnel run futurefunded-prod
user=elCUCO
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=20
stopsignal=TERM
stdout_logfile=/home/elCUCO/futurefunded/.logs/supervisor-cloudflared.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=/home/elCUCO/futurefunded/.logs/supervisor-cloudflared.err.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
